# weechat-relay/rootfs/etc/services.d/weechat/run
#!/command/with-contenv bash
set -e

CONFIG_DIR="/data/weechat"
OPTIONS_FILE="/data/options.json"

mkdir -p "${CONFIG_DIR}"

RELAY_PORT=$(jq -r '.relay_port' ${OPTIONS_FILE})
RELAY_PASSWORD=$(jq -r '.relay_password' ${OPTIONS_FILE})

echo "[INFO] WeeChat config dir: ${CONFIG_DIR}"
echo "[INFO] Relay port: ${RELAY_PORT}"

# Crear archivo de comandos para bootstrap
BOOTSTRAP_CMDS="${CONFIG_DIR}/bootstrap.weechat"

cat > "${BOOTSTRAP_CMDS}" <<EOF
/plugin load relay
/set irc.server_default.charset_decode iso-8859-1
/set irc.server_default.charset_encode iso-8859-1
/set relay.network.password ${RELAY_PASSWORD}
/set relay.network.bind_address 0.0.0.0
/set relay.network.port ${RELAY_PORT}
/relay add weechat 0.0.0.0/${RELAY_PORT}
/save
/quit
EOF

# Solo ejecutar bootstrap si no existe configuración previa
if [ ! -f "${CONFIG_DIR}/weechat.conf" ]; then
    echo "[INFO] First run detected — bootstrapping configuration"
    weechat -d "${CONFIG_DIR}" --no-terminal -r "/input eval ${BOOTSTRAP_CMDS}"
fi

echo "[INFO] Starting WeeChat daemon..."
exec weechat -d "${CONFIG_DIR}" --no-terminal
