# weechat-relay/rootfs/etc/services.d/weechat/run
#!/command/with-contenv bash
set -e

CONFIG_DIR="/data/weechat"
OPTIONS_FILE="/data/options.json"

mkdir -p "${CONFIG_DIR}"

RELAY_PORT=$(jq -r '.relay_port' "${OPTIONS_FILE}")
RELAY_PASSWORD=$(jq -r '.relay_password' "${OPTIONS_FILE}")

echo "[INFO] Using relay port: ${RELAY_PORT}"
echo "[INFO] WeeChat config dir: ${CONFIG_DIR}"

if [ ! -f "${CONFIG_DIR}/weechat.conf" ]; then
  echo "[INFO] First start: generating WeeChat config with ISO-8859-1 and relay"

  weechat -d "${CONFIG_DIR}" -r "
/set irc.server_default.charset_decode iso-8859-1
/set irc.server_default.charset_encode iso-8859-1
/set relay.network.password ${RELAY_PASSWORD}
/set relay.network.bind_address 0.0.0.0
/set relay.network.port ${RELAY_PORT}
/relay add weechat 0.0.0.0/${RELAY_PORT}
/save
/quit"
fi

export WEECHAT_HOME="${CONFIG_DIR}"

echo "[INFO] Starting WeeChat..."
exec weechat -d "${CONFIG_DIR}" --no-terminal
