# weechat-relay/rootfs/etc/services.d/weechat/run
#!/command/with-contenv bash
set -e

CONFIG_DIR="/data/weechat"

mkdir -p "${CONFIG_DIR}"

echo "[INFO] WeeChat config dir: ${CONFIG_DIR}"
echo "[INFO] Starting WeeChat (daemon, no terminal)..."

exec weechat -d "${CONFIG_DIR}" --no-terminalexec weechat -d "${CONFIG_DIR}" --no-terminal
