# weechat-relay/rootfs/etc/services.d/weechat/run
#!/command/with-contenv bash
set -ex

# Crear carpeta de configuración
mkdir -p /config/weechat

echo "[INFO] DEBUG MODE - El contenedor no se cerrará"

# Mantener contenedor vivo
tail -f /dev/null
