# weechat-relay/rootfs/etc/services.d/weechat/run
#!/command/with-contenv bash
set -e

CONFIG_DIR="/config/weechat"

mkdir -p "$CONFIG_DIR"
chmod -R 777 "$CONFIG_DIR"

echo "[INFO] Using config dir: $CONFIG_DIR"

# Primera inicialización: crear configuración mínima + desactivar plugins
if [ ! -f "$CONFIG_DIR/weechat.conf" ]; then
    echo "[INFO] Creating minimal config for first start"

    mkdir -p "$CONFIG_DIR"

    cat <<EOF > "$CONFIG_DIR/weechat.conf"
[plugin]
autoload = ""
EOF

    cat <<EOF > "$CONFIG_DIR/plugins.conf"
[desc]
autoload = "off"

[plug]
fifo = off
logger = off
irc = off
xfer = off
exec = off
charset = off
trigger = off
fset = off
spell = off
relay = on
EOF
fi

echo "[INFO] Starting WeeChat daemon with plugins disabled..."
exec weechat -d "$CONFIG_DIR" --no-terminal
